# desktop/source/bspwm/config/systemd-service/apt-update.service
[Unit]
Description=Run apt update

[Service]
Type=oneshot
ExecStart=/usr/bin/apt update

[Install]
WantedBy=multi-user.target
