# desktop/source/bspwm/config/systemd-service/apt-update.timer
[Unit]
Description=Run apt-update.service every 30 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
